# services/rc.local
#!/bin/sh -e
#
# Put this in /etc/rc.local
#

cd /var/prawnalith/local_images
sudo -u pi docker-compose up &
sleep 60
cd /var/prawnalith/services
sudo -u pi sh /var/prawnalith/services/start.sh &

exit 0

